# templates/default/statsd.conf.erb
description "StatsD"
author "Chef"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
expect daemon

exec start-stop-daemon --start --oknodo --chdir /usr/share/statsd --user statsd --chuid statsd --name statsd --startas /usr/local/bin/node --background -- /usr/share/statsd/stats.js /etc/statsd/config.js 2>&1 >> <%= @log_file %>

emits statsd-running
