# templates/default/statsd.conf.erb
description "StatsD"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

chdir /usr/share/statsd
<% unless @platform_version < 11.10 -%>
setuid statsd
<% end %>

script
<% if @platform_version < 11.10 -%>
  exec start-stop-daemon --start --chuid statsd --exec node /usr/share/statsd/stats.js /etc/statsd/config.js >> <%= @log_file %> 2>&1
<% else -%>
  exec node /usr/share/statsd/stats.js /etc/statsd/config.js >> <%= @log_file %> 2>&1
<% end -%>
end script

emits statsd-running
