# templates/default/statsd.conf.erb
description "StatsD"
author "Chef"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

chdir /usr/share/statsd
setuid statsd

script
  exec /usr/local/bin/node /usr/share/statsd/stats.js /etc/statsd/config.js 2>&1 >> <%= @log_file %>
end script

emits statsd-running
